# .config/autostart/configure-X-application-settings.desktop
[Desktop Entry]
Name=xrdb-runner
Comment=Run initial xrdb to set xterm dimensions
Type=Application
Exec=env bash -c 'xrdb -merge ~/.Xdefaults'
