# .config/autostart/remap-caps_lock-to-control.desktop
[Desktop Entry]
Name=Caps Lock remapper
Comment=Remaps Caps Lock to Control
Type=Application
Exec=setxkbmap -option "ctrl:nocaps"
