# .config/autostart/start-tmux.desktop
[Desktop Entry]
Name=tmux-starter
Comment=Start tmux on login
Type=Application
Exec=tmux start-server
